# examples/ubuntu-systemd/etc/profile.d/Z99-systemd.sh
#!/usr/bin/env bash

sudo /usr/local/bin/wsl-systemd-start